# nuki-bluetooth-watcher.service
[Unit]
Description=Nuki Bluetooth Watcher - Auto-restart Bluetooth when triggered
After=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/bin/nuki-bluetooth-watcher.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
